# bin/wmi.service
# 目录  /etc/systemd/system/
[Unit]
Description=Wmi Http Server
After=network.target
After=syslog.target

[Service]
Type=simple
User=www
Group=www
LimitNOFILE=65535
ExecStart=/usr/bin/php /opt/servers/echo/server.php
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target graphical.target